# peering/clodex.service
[Unit]
Description=Clodex (headless engine) — visual multi-agent PTY manager
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/wb-wrap-ui
# claude CLI lives in the user-global npm prefix; put it on PATH for spawns.
# node itself is resolved off this PATH by /usr/bin/env below (NodeSource drops
# it at /usr/bin/node; a nvm/local install stays reachable via ~/.npm-global etc).
Environment=PATH=%h/.npm-global/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
# The remote node runs the HEADLESS engine (plain Node, no Electron/Xvfb): a peer
# node needs only the engine + node-pty, never the desktop GUI. headless-main.js
# stands up the same engine the Electron app hosts and answers the peer wire.
ExecStart=/usr/bin/env node headless-main.js
Restart=always
RestartSec=5
# Give the app time to release its single-instance lock + unix sockets on stop
# before systemd would consider a restart, mirroring the seed script's sleep.
TimeoutStopSec=15

[Install]
WantedBy=default.target
